# supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/run-mysqld

[program:cron]
command=/usr/sbin/cron -f