# supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/run-mysqld
redirect_stderr=true
stdout_logfile=/mysql/mysql.err
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
autorestart=true

[program:cron]
command=/usr/sbin/cron -f
